# examples/systemd/qdrant.service
# /etc/systemd/system/qdrant.service
#
# Single-node Qdrant for an agent setup. Storage under /var/lib/qdrant.
# Bind 127.0.0.1 only — never expose Qdrant to the public internet.
#
# Install:
#   sudo cp examples/systemd/qdrant.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now qdrant
#
# Verify:
#   curl -sS http://127.0.0.1:6333/healthz   # → "healthz check passed"

[Unit]
Description=Qdrant Vector Database
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qdrant
Group=qdrant
WorkingDirectory=/var/lib/qdrant
Environment=QDRANT__SERVICE__HTTP_PORT=6333
Environment=QDRANT__SERVICE__GRPC_PORT=6334
Environment=QDRANT__SERVICE__HOST=127.0.0.1
Environment=QDRANT__STORAGE__STORAGE_PATH=/var/lib/qdrant/storage
ExecStart=/usr/local/bin/qdrant
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
